# src/etc/services.d/writable-paths/up
with-contenv -ew
multisubstitute {
    importas -D "" UNSPLITTED WRITABLE_PATHS
    importas -sd":" PATHS WRITABLE_PATHS
    importas -D "" USER WRITABLE_USER
}
if -t {
    s6-test -n $UNSPLITTED -a -n $USER
}
if {
    s6-notice "Making ${UNSPLITTED} writable for ${USER}"
    forx -p TARGET { $PATHS }
        importas -u TARGET TARGET
        if {
            if -t {
                s6-test ! -e $TARGET
            }
            s6-logcmd
            s6-mkdir -p $TARGET
        }
        ifelse {
            s6-test -d $TARGET
        } {
            s6-logcmd
            setfacl -R -m u:${USER}:rwX -m d:u:${USER}:rwX -m mask:rwX -m d:mask:rwX ${TARGET}
        }
        s6-logcmd
        setfacl -m u:${USER}:rwX -m mask:rwX ${TARGET}
}
importas -i LAUNCH_ENV S6_LAUNCH_ENV
if {
    s6-rmrf ${LAUNCH_ENV}/env/WRITABLE_PATHS
}
s6-rmrf ${LAUNCH_ENV}/env/WRITABLE_USER
