# src/etc/services.d/writable-paths/up
with-contenv -ew
multisubstitute {
    importas -D "" UNSPLITTED WRITABLE_PATHS
    import -sd":" WRITABLE_PATHS
    import -D "" WRITABLE_USER
}
if -t {
    s6-test -n $UNSPLITTED -a -n $WRITABLE_USER
}
if {
    s6-notice "Making ${UNSPLITTED} writable for ${WRITABLE_USER}"
    forx -p TARGET { $WRITABLE_PATHS }
        import -u TARGET
        if {
            if -t {
                s6-test ! -e $TARGET
            }
            s6-logcmd
            s6-mkdir -p $TARGET
        }
        ifelse {
            s6-test -d $TARGET
        } {
            s6-logcmd
            setfacl -R -m u:${WRITABLE_USER}:rwX -m d:u:${WRITABLE_USER}:rwX ${TARGET}
        }
        s6-logcmd
        setfacl -m u:${WRITABLE_USER}:rwX ${TARGET}
}
import -i S6_LAUNCH_ENV
if {
    s6-rmrf ${S6_LAUNCH_ENV}/env/WRITABLE_PATHS
}
s6-rmrf ${S6_LAUNCH_ENV}/env/WRITABLE_USER
