# overlay/etc/services.d/entrypoint.d/up
import-launch-var SKIP_ENTRYPOINT_D
multisubstitute {
    importas -i DIR S6_ENTRYPOINT_D
    importas -D "" SKIP SKIP_ENTRYPOINT_D
}
if -t { s6-test -z $SKIP -a -d $DIR }
s6-progress "Running ${DIR} commands"
forbacktickx -0 __SCRIPT { pipeline { s6-ls -0 $DIR } s6-sort }
importas -u SCRIPT __SCRIPT
with-contenv -a
if -t { s6-test -x ${DIR}/${SCRIPT} }
s6-logcmd
${DIR}/${SCRIPT}
